# systemd.pixraper@.service
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pixraper.glue --add --devname=/dev/%I
